# deploy/supervisor/xml-mapper-worker.conf
[program:xml-mapper-worker]
command=/usr/bin/php /var/www/xml-mapper/artisan queue:work redis --queue=imports,normalization,feeds,dictionaries --sleep=3 --tries=3 --timeout=1800 --max-time=3600
directory=/var/www/xml-mapper
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/xml-mapper/worker.log
stopwaitsecs=3600
