# deploy/cron/schedule-run.cron
* * * * * cd /var/www/xml-mapper/current && php artisan schedule:run >> /var/log/xml-mapper/scheduler.log 2>&1
